# deploy/services/homer.service
[Unit]
Description=coturn_docker
After=docker.service
Requires=docker.service
After=network-online.target
Requires=network-online.target

[Service]
Type=simple
Restart=always
TimeoutStartSec=0
RestartSec=1
WorkingDirectory=/sipmediagw/deploy
EnvironmentFile=/etc/environment
ExecStart=/usr/bin/docker compose -p homer up --force-recreate heplify_server homer_webapp
ExecStop=/usr/bin/docker compose -p homer down

[Install]
WantedBy=multi-user.target
